# rootfs/etc/s6/sshd/run
#!/bin/bash

# Don't detach (-D), log to stderr (-e)
exec /usr/sbin/sshd -D -e
